# systemd/pacman-offline-prepare.timer
[Unit]
Description=Prepare pacman offline system-update

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
