# systemd/pacman-offline-prepare.service
[Unit]
Description=Prepare pacman offline system-update
ConditionPathExists=!/var/lib/pacman/db.lck

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman-offline -y
KillSignal=SIGINT
ProtectSystem=on
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on
